# packages/graphics/libmali/system.d/libmali-config.service
[Unit]
Description=Runtime symlink service for Mali library
Before=graphical.target
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/libmali-config'

[Install]
WantedBy=multi-user.target
